# systemd/.config/systemd/user/mbsync.path
[Path]
PathChanged=%h/.local/share/mail/.notmuch

[Install]
WantedBy=default.target
